# base/docker/conf/supervisord-du.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:netopeer2-server]
command=/usr/local/sbin/netopeer2-server -d -v 3 -t 60 -x /data/ietf-yang-schema-mount.xml
stdout_logfile=/var/log/netopeer2-server.log
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/sshd.log
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0
autorestart=false


[program:vsftpd]
command=/usr/sbin/vsftpd
stdout_logfile=/var/log/vsftpd.log
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0
autorestart=true
